# supervisord.conf
[supervisord]
nodaemon=true

[program:haproxy]
command=/usr/sbin/haproxy -f /etc/haproxy/haproxy.cfg

[program:chat]
command=/usr/local/bin/node /opt/app/index.js
